# roles/nextcloud/templates/nextcloud_cron.timer
[Unit]
Description=Nextcloud cron job timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
